# configs/service/supervisor-app.conf
[program:service]
priority=5
command = python3 service/run.py