# configs/django_app/supervisor-app.conf
[program:uwsgi]
priority = 4
command = uwsgi --ini /home/docker/code/uwsgi.ini
